# auroraos/config/includes.chroot/etc/xdg/autostart/aurora-torbrowser.desktop
[Desktop Entry]
Type=Application
Name=AuroraOS Tor connection assistant (Tor mode)
Comment=In Tor mode, connect to Tor (direct or via a bridge); opens Tor Browser once connected
Exec=/usr/local/bin/aurora-tor-assistant
Terminal=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
# Short delay before the assistant script starts. The script itself then waits
# for the GNOME Initial Setup wizard (language -> keyboard -> Wi-Fi) to finish
# before showing anything, so the "connect to Tor" step reads as the NEXT step
# after the wizard rather than popping up on top of it.
X-GNOME-Autostart-Delay=2
